# scripts/services/traewelling-scheduler.timer
[Unit]
Description=Traewelling Scheduler

[Timer]
OnCalendar=minutely

[Install]
WantedBy=timers.target
